# src/slash-bedrock/strata/fallback/etc/inittab
### Bedrock Linux 1.0beta2 Nyla
### /etc/inittab

# run on boot just after the kernel finishes loading
::sysinit:/etc/init.d/rcS

# ttys
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4
tty5::respawn:/sbin/getty 38400 tty5
tty6::respawn:/sbin/getty 38400 tty6

# run on reboot
::ctrlaltdel:/bin/kill 1
# run just before shutdown
::shutdown:/etc/init.d/rcK
